# tools/hotplug/Linux/systemd/xenpcid.service.in
[Unit]
Description=Xen PCI host daemon
ConditionVirtualization=xen

[Service]
Type=forking
ExecStart=@sbindir@/xl pcid

[Install]
WantedBy=multi-user.target
